# external/rkscript/61-sd-cards-auto-mount.rules
KERNEL!="mmcblk*[0-9]", GOTO="sd_cards_auto_mount_end"
SUBSYSTEM!="block", GOTO="sd_cards_auto_mount_end"
ACTION=="add", PROGRAM!="/sbin/blkid %N", GOTO="sd_cards_auto_mount_end"
ATTRS{type}!="SD", GOTO="sd_cards_auto_mount_end"

IMPORT{program}="/sbin/blkid -o udev -p %N"

ACTION=="add", ENV{ID_FS_TYPE}=="vfat",  RUN+="/usr/bin/mount.sdcard.sh vfat /dev/%k"
ACTION=="add", ENV{ID_FS_TYPE}=="ntfs",  RUN+="/usr/bin/mount.sdcard.sh ntfs /dev/%k"
ACTION=="add", ENV{ID_FS_TYPE}=="exfat", RUN+="/usr/bin/mount.sdcard.sh exfat /dev/%k"
ACTION=="add", ENV{ID_FS_TYPE}=="ext2",  RUN+="/usr/bin/mount.sdcard.sh ext2  /dev/%k"
ACTION=="add", ENV{ID_FS_TYPE}=="ext3",  RUN+="/usr/bin/mount.sdcard.sh ext3  /dev/%k"
ACTION=="add", ENV{ID_FS_TYPE}=="ext4",  RUN+="/usr/bin/mount.sdcard.sh ext4  /dev/%k"

ACTION=="remove", RUN+="/bin/umount '/mnt/sdcard'"

LABEL="sd_cards_auto_mount_end"
